# geloge-lo-crawl/service/geloge-lo-crawl-tweet-deliver/run
#!/bin/bash
GELOGE_LO_CRAWL=../../geloge-lo-crawl-tweet-deliver
exec 2>&1
exec $GELOGE_LO_CRAWL
